# sysadmin/root.crontab
#!/bin/sh
# Avoid editing this directly on the server. Instead, edit this and then
# run `sysadmin/rebuild.sh` which will update the crontab on the server
22 3 * * * bash /home/web/app/crawler/RefreshCourses.sh
43 2 * * * bash /home/web/app/sysadmin/phpmyadmin/update.sh
51 4 * * * letsencrypt renew && systemctl reload nginx
59 6 * * * bash /home/web/app/sysadmin/backup-upload.sh
